# systemd/bravo-coder.service
[Unit]
Description=Bravo Pipeline — Coder Agent
Documentation=https://github.com/ppxtestarena4/pipeline-template
After=network-online.target
Wants=network-online.target
# Запускаться после сетевого соединения

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/home/agent/<repo>

# Путь к скрипту агента
ExecStart=/home/agent/<repo>/pipeline/coder-daemon.sh

# Переменные окружения
Environment=OPENAI_API_KEY=
Environment=PIPELINE_REPO=ppxtestarena4/pipeline-template
Environment=SLEEP_INTERVAL=300
Environment=HOME=/home/agent
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/agent/.local/bin

# Политика перезапуска
Restart=always
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

# Вывод логов в systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bravo-coder

# Таймаут ожидания при остановке
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
